# deploy/horizon.supervisor.conf
; LEXA Horizon worker, supervised so it stays up across crashes.
;
; Place at /etc/supervisord.d/lexa-horizon.ini (provision.sh does this).
; The cPanel user must be writable to /home/__CPANEL_USER__/lexa.

[program:lexa-horizon]
process_name=%(program_name)s_%(process_num)02d
directory=/home/__CPANEL_USER__/lexa
command=/usr/bin/env php artisan horizon
autostart=true
autorestart=true
user=__CPANEL_USER__
numprocs=1
redirect_stderr=true
stdout_logfile=/home/__CPANEL_USER__/lexa/storage/logs/horizon.log
stopwaitsecs=3600
environment=APP_ENV="production"
